# supervisord.conf
; supervisord — orquestra os 2 serviços do auton-painel (PID 1 do container).
; Cada serviço tem seu próprio CWD → cada um lê seu próprio /app/<dir>/.env
; (montado via compose). Logs vão pro stdout/stderr → Docker captura.

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[program:concorrentes]
directory=/app/concorrentes
command=node .mastra/output/index.mjs
autorestart=true
startsecs=10
startretries=3
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PORT=4111,NODE_ENV=production

[program:batalhao]
directory=/app/batalhao
command=npx tsx src/server.ts
autorestart=true
startsecs=10
startretries=3
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PORT_BATALHAO=4112,NODE_ENV=production
